# adgs/common_ADGS/ADGSConf/confPdsImport/pds2_import_pds-imp1-5.timer
[Unit]
Description=Run PDS2 Import PDS-IMP1-5 every 1 minute

[Timer]
OnCalendar=minutely
Persistent=true
Unit=pds2_import_pds-imp1-5.service

[Install]
WantedBy=timers.target
