# adgs/common_ADGS/initScripts/systemd/adgs-catalogueexportgenerator.timer
[Unit]
Description=ADGS Catalog Exporter Generator
Documentation=https://www.exprivia.com
Requires=network.target
After=network.target

[Timer]
OnCalendar=*-*-* 03:00:00
Unit=adgs-catalogueexportgenerator.service

[Install]
WantedBy=timers.target

