# auxip/auxip-subscription-backend/setup/rpm/extrafiles.list
%{INST_DIR}/%{JAR_FILENAME}

/usr/lib/systemd/system/%{SERVICE_FILENAME}
