# adgs/applications/DataRetriever/setup/conf/supervisor.d/data-retriever-download-worker.ini
[program:download-worker-{project}]
environment=PHP_IMPORT_CONFIGURATION="{config}"
command=php {installpath}/{installname}/workers/data-retriever-download-worker.php
process_name=%(program_name)s-%(process_num)d
numprocs={downloadworkers}
directory={installpath}/{installname}/workers
user={daemonuser}
autorestart=true
startretries=10000
redirect_stderr=true
stdout_syslog = true
stdout_logfile = /var/log/supervisor/data-retriever-download-worker-{project}.log
loglevel=info
stdout_logfile_maxbytes=1MB
